# etc/script/systemd/baremetal-api.service
[Unit]
Description=baremetal api service

[Service]
Restart=always
User=root
ExecStart=/usr/bin/baremetal-api

[Install]
WantedBy=multi-user.target